# debian/systemd/browsermcp-daemon.service
[Unit]
Description=BrowserMCP Enhanced - WebSocket Daemon
After=network-online.target browsermcp-http.service
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__GROUP__
Environment=NODE_ENV=production
Environment="HOME=__HOME__"
EnvironmentFile=/etc/default/browsermcp
WorkingDirectory=__INSTALL_DIR__
ExecStart=/usr/bin/node __INSTALL_DIR__/dist/daemon/websocket-daemon.js
Restart=always
RestartSec=3
KillMode=process

[Install]
WantedBy=multi-user.target

